# src/hwaetd.rc.sh
#!/bin/ksh

daemon="/usr/local/sbin/hwaetd"

. /etc/rc.d/rc.subr

rc_cmd $1
